# root/etc/s6-overlay/s6-rc.d/svc-grok2api-log/run
#!/usr/bin/with-contenv bash

exec logutil-service /app/logs
